# systemd/ai-update.timer
[Unit]
Description=AI Update - Daily at 3:00 AM (NOT hourly)

[Timer]
# 每天凌晨 3:00 运行一次（不是每小时！）
OnCalendar=daily
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
